# diviner-live.service
[Unit]
Description=Diviner live trading bot
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/trading
EnvironmentFile=/home/admin/trading/trader.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/env node src/live.js --strategy src/strategies/release_morning_short_testing.js --ticker SBER --sandbox --account CHANGE_ME_ACCOUNT_ID --log /home/admin/trading/trader.log
Restart=always
RestartSec=10
KillSignal=SIGTERM
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
